# .ddev/web-build/email.cron
* * * * * IS_DDEV_PROJECT=true DDEV_PHP_VERSION=8.3 /var/www/html/bin/cron
